# aws/systemd/drift-base.service
[Unit]
Description=Drift Base Server
After=syslog.target

[Service]
Type=notify
EnvironmentFile=/etc/environment
ExecStart=/usr/local/bin/uwsgi /etc/opt/drift-base/config/uwsgi.ini
Restart=on-failure
RestartSec=5
KillSignal=SIGQUIT
User=root
NotifyAccess=all

[Install]
WantedBy=multiuser.target
